# crontab
12 0 * * * python2.7 ~/webapps/ucl/uchicagolunch/manage.py prunequeue
@daily python2.7 ~/webapps/ucl/uchicagolunch/manage.py match
@daily python2.7 ~/webapps/ucl/uchicagolunch/manage.py loadweek
